# root/etc/services.d/gitwatch/run
#!/usr/bin/execlineb -P

with-contenv
s6-setuidgid abc
/usr/local/bin/run_gitwatch
